# third_party/service/ifconfig_setup.service
[Unit]
Description=Check and update network configuration

[Service]
Type=simple
ExecStart=/usr/bin/ifconfig_setup.sh

[Install]
WantedBy=default.target
